# docker_cron/crontab
1 */1 * * * /bin/sh /root/scripts/backup.sh

